# i-spy-notify.service
[Unit]
Description=I Spy Notify notifications observer
Documentation=man:i-spy-notify(1)

[Service]
Type=simple
ExecStart=/usr/local/bin/i-spy-notify
Restart=on-failure

[Install]
WantedBy=default.target
